# scripts/com.bluevlad.auto-tobe-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Auto-Tobe-Agent launchd 기본 설정 (레거시 호환)
  이슈 수정 배치를 하루 3회 실행합니다.

  설치:
    1. 이 파일의 경로를 운영 환경에 맞게 수정하세요 (AGENT_HOME, HOME 등)
    2. cp scripts/com.bluevlad.auto-tobe-agent.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.bluevlad.auto-tobe-agent.plist

  주의: 경로에 실제 사용자명/서버 정보를 포함하지 않도록 local.json에서 오버라이드하세요.
-->
<dict>
    <key>Label</key>
    <string>com.bluevlad.auto-tobe-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>${AGENT_HOME}/scripts/run-batch.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>${AGENT_HOME}</string>

    <key>StandardOutPath</key>
    <string>${AGENT_HOME}/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>${AGENT_HOME}/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>AbandonProcessGroup</key>
    <false/>
</dict>
</plist>
